# os/modules/shell-init/src/usr/lib/systemd/user/sideral-shell-seed.service
[Unit]
Description=sideral shell profile seeding
Documentation=https://github.com/athenabriana/sideral
# Run after the basic session environment is up but before the user's
# shell or graphical session starts consuming the config files.
After=default.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/sideral-shell-seed
# Non-fatal: if the script errors (e.g. usermod requires sudo), the
# session still starts. The seed is best-effort on each login.
RemainAfterExit=yes

[Install]
WantedBy=default.target
