# system_files/usr/lib/systemd/user/athens-vscode-setup.service
[Unit]
Description=Install VS Code extensions on first user login
After=graphical-session.target
# Skip if the marker already exists for this user.
ConditionPathExists=!%h/.cache/athens/vscode-setup-done

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c '\
    mkdir -p %h/.cache/athens && \
    command -v code >/dev/null 2>&1 || exit 0 && \
    code --install-extension ms-vscode-remote.remote-ssh --force && \
    code --install-extension ms-vscode-remote.remote-containers --force && \
    code --install-extension ms-azuretools.vscode-containers --force && \
    touch %h/.cache/athens/vscode-setup-done'
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
