# system_files/usr/lib/systemd/user/athens-mise-install.service
[Unit]
Description=Install mise + core CLI tools on first user login
Documentation=https://mise.jdx.dev
After=graphical-session.target
# Marker keeps this unit idempotent — runs once, retries if it previously failed
# (e.g., offline at first login).
ConditionPathExists=!%h/.cache/athens/mise-setup-done

[Service]
Type=oneshot
Environment=MISE_QUIET=1
ExecStart=/usr/bin/bash -c '\
    set -e; \
    mkdir -p %h/.cache/athens; \
    if [ ! -x %h/.local/bin/mise ]; then \
        curl -sSfL https://mise.run | sh; \
    fi; \
    export PATH="%h/.local/bin:$PATH"; \
    mise install act atuin direnv; \
    touch %h/.cache/athens/mise-setup-done'
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
